# bin/supervisord.conf
[supervisord]
nodaemon=true

[program:gerrit]
command=/bin/bash -c "/home/gerrit/site/bin/gerrit.sh start"
